# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:api]
command=/usr/bin/dotnet /app/api/Comprendo.Api.dll
environment=ASPNETCORE_ENVIRONMENT="Production",ASPNETCORE_URLS="http://127.0.0.1:8080"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=10

[program:bot]
command=/usr/bin/node /app/bot/src/index.js
directory=/app/bot
environment=PORT="3001",NODE_ENV="production",CORE_API_URL="http://127.0.0.1:8080",INTEGRATION_API_KEY="%(ENV_INTEGRATION_API_KEY)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",ZHIPU_API_KEY="%(ENV_ZHIPU_API_KEY)s",GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",GEMINI_API_KEY="%(ENV_GEMINI_API_KEY)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=20

[program:frontend]
command=/usr/bin/node /app/frontend/server.js
directory=/app/frontend
environment=PORT="3000",HOSTNAME="127.0.0.1",NODE_ENV="production"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=30

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=40
